# rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bash
#shellcheck shell=bash

exec \
  /usr/sbin/nginx 2>&1 | awk -W interactive '{print "[nginx] " $0}'
